# packages/httpdatetime/httpdatetime-1.0.1-py3-none-any.whl/httpdatetime/statics/httpdatetime.service
[Unit]
Description=Launch httpdatetime to synchronize date and time online using HTTPS.

[Service]
Type=simple
ExecStart=httpdatetime

[Install]
WantedBy=multi-user.target
